# src/debian_control/prerm
#!/bin/sh
set +e

service "LiberaSinglePass" stop

if [ "$1" = "remove" ]; then
    update-rc.d -f "LiberaSinglePass" remove

    rm /etc/init.d/LiberaSinglePass
fi
